# slackbot.service
[Unit]
Description=Slackbot

[Service]
ExecStart=/home/morty/errbot-venv/errbot-root/start_bot

[Install]
WantedBy=multi-user.target
